# supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=/usr/sbin/php-fpm -F

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
